# templates/purge_sssd.service.erb
[Unit]
Description=Purges SSSD cache on boot prior to autofs starting
Before=multi-user.target sssd.service autofs.service systemd-user-sessions.service nss-user-lookup.target
Wants=nss-user-lookup.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/root/purge_sssd start
TimeoutSec=180

[Install]
WantedBy=multi-user.target
